# templates/default/sysvinit.erb
#!/bin/bash
#
# OPT
# Managed by Chef for node <%= node['fqdn'] -%>.

# chkconfig: 2345 82 55
# processname: <%= @config.name %>
# description: Runs puma for rails.
NAME=<%= @config.name %>
. /etc/rc.d/init.d/functions
RAILS_ENV=<%= @config.rails_env %>
APP_PATH=<%= "#{@config.deploy_path}/current" %>
PID_PATH=/var/opt/shared

PID=$PID_PATH/pids/$NAME.pid
LOCK=/var/lock/subsys/$NAME


start() {
  cd $APP_PATH

  # Start puma
  echo -n $"Starting $NAME: "
#  daemon "PATH=<%= @config.bundler_path %>:$PATH RAILS_ENV=$RAILS_ENV bundle exec thin start -d -u <%= @config.run_user %> -g <%= @config.run_group %> -p <%= @config.port %> -P $PID"
  daemon "PATH=<%= @config.bundler_path %>:$PATH RAILS_ENV=$RAILS_ENV bundle exec puma -d -w 2 -p <%= @config.port %> --pidfile $PID"
  puma=$?
  [ $puma -eq 0 ] && touch $LOCK
  echo
  retval=$puma
  return $retval
}

stop() {
  cd $APP_PATH

  # Stop puma
  echo -n $"Stopping $NAME: "
  killproc -p $PID
  puma=$?
  [ $puma -eq 0 ] && rm -f $LOCK
  echo
  # return retval
  return $puma
}

restart() {
  stop
  start
}

get_status() {
  status -p $PID $NAME
}

query_status() {
  get_status >/dev/null 2>&1
}

case "$1" in
  start)
    query_status && exit 0
    start
    ;;
  stop)
    query_status || exit 0
    stop
    ;;
  restart)
    restart
    ;;
  status)
    get_status
    ;;
  *)
    N=/etc/init.d/$NAME
    echo "Usage: $N {start|stop|restart|status}" >&2
    exit 1
    ;;
esac

exit 0
